# files/xray-tls-certs.timer
[Unit]
Description=Schedules xray-tls-certs.service

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
